# pi-dispenser/systemd/pi-dispenser.service
# HealMeDaily pill-dispenser agent (Phase 8)
# Install:
#   sudo cp systemd/pi-dispenser.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now pi-dispenser
# Logs: journalctl -u pi-dispenser -f

[Unit]
Description=HealMeDaily pill dispenser agent (not a certified medical device)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-dispenser
# .env holds DISPENSER_MEDPLUM_BASE_URL / _CLIENT_ID / _CLIENT_SECRET / _PATIENT_ID
EnvironmentFile=/home/pi/pi-dispenser/.env
ExecStart=/home/pi/pi-dispenser/.venv/bin/python -m pi_dispenser run
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
